# systemd/karaoke420.service
[Unit]
Description=Karaoke 420
After=remote-fs.target network-online.target
Wants = network-online.target

[Service]
User=osmc
Group=osmc
Type=simple
WorkingDirectory=/home/osmc/karaoke420
ExecStart=/bin/bash -lc 'rails server -b 0.0.0.0'

[Install]
WantedBy=multi-user.target
